# debian/tuxedo-daemon.service
[Unit]
Description=TUXEDO Hardware Control Daemon
After=multi-user.target
Documentation=man:tuxedo-daemon(8)

[Service]
Type=dbus
BusName=com.tuxedo.Control
ExecStart=/usr/bin/tuxedo-daemon
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/sys/devices/platform/tuxedo_io /sys/devices/system/cpu /sys/class/backlight /sys/class/leds
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
